# upstart/wego.conf
description "Hope word filter service"

start on runlevel [2345]
stop on runlevel [016]

respawn
respawn limit 3 30

exec /home/test/wego